# dist/upstart/quantum.conf
# quantum -- quantum sdn service
#
# looks for /etc/quantum/quantum.env or /etc/default/quantum to configure the sdn

description "quantum sdn service"

start on ( started networking )
stop on ( stopped networking )

kill timeout 20

respawn

env PID_FILE="/var/run/quantum.pid"
env QUANTUM="/usr/sbin/quantum"

script
    trap "kill -HUP $(cat ${PID_FILE})" HUP

    set -a

    if [ -f /etc/quantu/quantum.env ]; then
        . /etc/quantum/quantum.env
    fi

    if [ -f /etc/default/quantum ]; then
        . /etc/default/quantum
    fi

    exec start-stop-daemon --pidfile ${PID_FILE} --exec ${QUANTUM} --start -- ${QUANTUM_CLI_OPTS}
end script

post-stop script
    rm -f ${PID_FILE}
end script
